# startrobotino.conf
#/usr/share/upstart/sessions
description "my script"
start on started unity8-dash 

exec /home/phablet/run-pocketsphinx.sh
